# launchd/com.vokomarov.cc-approvals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vokomarov.cc-approvals</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vokomarov/go/bin/cc</string>
        <string>daemon</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/cc-approvals.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cc-approvals-error.log</string>
</dict>
</plist>
