# launchd/com.vokomarov.cc-approvals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vokomarov.cc-approvals</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vovan/go/bin/claude-code-approvals</string>
        <string>daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vovan/go/src/github.com/vokomarov/claude-code-approvals</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/vovan/go/src/github.com/vokomarov/claude-code-approvals/logs/cc-approvals.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vovan/go/src/github.com/vokomarov/claude-code-approvals/logs/cc-approvals-error.log</string>
</dict>
</plist>
